# dockers/docker-dummy-http/dummy_http.conf
[program:dummy-http]
command=hello-sonic-service
autostart=true
autorestart=true
priority=10