# dockers/docker-dummy-http/dummy_http.conf
[program:dummy-http]
command=hello-sonic-service
autostart=true
autorestart=true
priority=10
stdout_logfile=NONE
stdout_syslog=true
stderr_logfile=NONE
stderr_syslog=true